# supervisord.conf
[supervisord]
nodaemon=true

[program:panaf]
redirect_stderr=true
command=apachectl restart

[program:cron]
redirect_stderr=true
command=service cron restart

